# etc/watchdog/supervisor_watchdog.conf
[program:pong]
command=sh /home/pong/src/circpong/etc/watchdog_client.sh
process_name=%(program_name)s
numprocs=1
directory=/tmp
autorestart=true
startsecs=2
startretries=1000000

